# uwsgi_start.ini
[program:uwsgi_start]

directory = /home/workspace/python_real_war
command = /root/.virtualenvs/realwar/bin/uwsgi --ini /home/workspace/python_real_war/uwsgi.ini

user = root


stdout_logfile = /home/workspace/python_real_war/log/uwsgi_out.log
stderr_logfile = /home/workspace/python_real_war/log/uwsgi_err.log
autostart = true
autorestart = true
startsecs = 0
stopwaitsecs = 0
stopsignal = QUIT
redirect_stderr = true
loglevel = error
logfile_maxbytes = 1M
